# tests/root/etc/crontab
MAILTO=cron@example.com
42  * * * * lutter /usr/local/bin/backup
54 16 * * * lutter /usr/sbin/stuff
